# server/packaging/shallan.service
[Unit]
Description=Shallan music server

[Service]
Type=exec
ExecStart=shalland
Restart=always

[Install]
WantedBy=multi-user.target
